# docs/demo/reset-demo.timer
[Unit]
Description=reset demo nightly

[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=timers.target
